# plex-streamer/crontab
35 12,19 * * * bash /scripts/plex/save_library_to_disk.sh
35 5 * * * bash /scripts/plex/save_library_to_disk.sh && bash /scripts/plex/load_master_library_db.sh