# plex-streamer/crontab
35 */8 * * * bash /scripts/plex/save_library_to_disk.sh